# chat.service
[Unit]
Description=Chat Platform Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/chat
ExecStart=/opt/chat/chat-server
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables (or use .env file in WorkingDirectory)
# Environment="PORT=8080"
# Environment="JWT_SECRET=your_secret_here"

[Install]
WantedBy=multi-user.target
